# logging/crontab
# restart goaccess after a new log file has been started
1 0 * * * root killall -TERM goaccess
